# givenergy_local/rootfs/etc/services.d/givenergy_local/run
#!/usr/bin/with-contenv bashio
LIVE_REFRESH=$(bashio::config 'live_refresh_seconds')
FULL_REFRESH=$(bashio::config 'full_refresh_seconds')
export GIVLOCAL_LIVE_REFRESH_SECONDS="${LIVE_REFRESH:-30}"
export GIVLOCAL_FULL_REFRESH_SECONDS="${FULL_REFRESH:-120}"
export GIVLOCAL_DATA_DIR="/data"

cd /app/backend
exec python3 -m uvicorn main:app --host 0.0.0.0 --port 8377
